# cron/crontab
0 * * * * php /src/cron/update_meters.php && python /src/main.py
